# Coder/deploy/coder.service
[Unit]
Description=Coder Agent FastAPI Server
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/coder-agent
EnvironmentFile=/opt/coder-agent/.env
ExecStart=/opt/coder-agent/.venv/bin/uvicorn Coder.server.main:app --host 127.0.0.1 --port 8000 --workers 2
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
